# .docker/app/supervisor/laravel-scheduler.conf
[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=/bin/sh -c "while [ true ]; do (php /app/backend/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
autostart=true
autorestart=true
numprocs=1
user=%(ENV_MY_USER)s
redirect_stderr=true
stdout_logfile=/app/backend/storage/logs/scheduler.log
